# systemd/w3p-ups.service
[Unit]
Description=Web3 Pi UPS Monitoring Service
Documentation=https://github.com/Web3-Pi/Web3-Pi-UPS-Service
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/w3p-ups --config /etc/w3p-ups/config.toml
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

# Runtime directory for IPC socket
RuntimeDirectory=w3p-ups
RuntimeDirectoryMode=0755

# Security hardening
NoNewPrivileges=false
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /var/run/w3p-ups
PrivateTmp=true

# Allow access to serial devices
SupplementaryGroups=dialout

[Install]
WantedBy=multi-user.target
